# systemd/user/phase1-daily-use.service.example
[Unit]
Description=Northstar Research Phase 1 Daily Research Tools Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/northstar-research/apps/phase1-daily-use
ExecStart=%h/.local/bin/uvicorn app:app --host 127.0.0.1 --port 3020
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
