# systemd/user/phase1-daily-use.service
[Unit]
Description=Northstar Research — Phase 1 daily services
After=network-online.target

[Service]
Type=oneshot
ExecStartPre=/bin/bash -c 'echo "Starting Northstar Research Phase 1 services"'
ExecStart=docker compose -f %h/northstar-research/docker/docker-compose.yml up -d
RemainAfterExit=yes

[Install]
WantedBy=default.target
