# systemd/user/research-portal-native.service
[Unit]
Description=Northstar Research Console (React PWA + API proxy)
After=network-online.target northstar-pg.service

[Service]
Type=simple
WorkingDirectory=%h/northstar-research
EnvironmentFile=%h/northstar-research/apps/research-portal/.env
ExecStart=%h/northstar-research/.venv/bin/uvicorn research_portal.main:app --host 0.0.0.0 --port 3010
Restart=on-failure

[Install]
WantedBy=default.target
