# systemd/user/research-portal-native.service.example
[Unit]
Description=Northstar Research Portal Native WSL Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/northstar-research/apps/research-portal
Environment=RESEARCH_API=http://127.0.0.1:8099
ExecStart=%h/.local/bin/uvicorn app:app --host 127.0.0.1 --port 3010
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
